# systemd/.config/systemd/user/dwall-watcher.path
[Path]
PathModified=/home/dylan/.local/share/dwall/style

[Install]
WantedBy=multi-user.target
