# systemd/.config/systemd/user/wttr-in.timer
[Unit]
Wants=graphical.target

[Timer]
OnBootSec=0s
OnUnitInactiveSec=15m

[Install]
WantedBy=timers.target
